# nethserver_fwrules.cron
# Cron file to reset firewall rules

0 0 * * * root /usr/share/cockpit/nethserver-fwrules-update-status/reset_rules.sh
